# agents/agents.service
# AgentOS - Agent Runner Systemd Service
# Runs the agent runner API on port 3001, auto-starts on boot.
#
# Install:
#   sudo cp agents.service /etc/systemd/system/agentos-agents.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now agentos-agents.service

[Unit]
Description=AgentOS Agent Runner
Documentation=https://github.com/Suneel-DK/AgentOS
After=network.target ollama.service agentos-dashboard.service
Wants=network.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/opt/agentos/agents
ExecStart=/usr/bin/python3 /opt/agentos/agents/agent-runner.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentos-agents

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
